# packaging/systemd/jacques.service
[Unit]
Description=Jacques automatic disc ripping daemon
After=network.target systemd-udev-settle.service

[Service]
Type=simple
User=jacques
Group=jacques
SupplementaryGroups=cdrom

EnvironmentFile=-/etc/jacques/jacques.conf
Environment=JACQUES_DB_PATH=/var/lib/jacques/jacques.db
Environment=JACQUES_OUTPUT_PATH=/var/lib/jacques/library
Environment=JACQUES_TEMP_PATH=/var/lib/jacques/tmp

ExecStartPre=/usr/bin/mkdir -p $JACQUES_OUTPUT_PATH $JACQUES_TEMP_PATH
ExecStart=/usr/bin/jacques

Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/jacques

[Install]
WantedBy=multi-user.target
